# init/snapr-snap.service
[Unit]
Description=Take snapshots with snapr

[Service]
Type=oneshot
ExecStart=/sbin/snapr --snap

[Install]
WantedBy=multi-user.target