# init/snapr-snap.timer
[Unit]
Description=Trigger hourly snapr snapshots

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target